# host/systemd/pachislot-data-counter.service
[Unit]
Description=Pachislot Data Counter (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nakano
Group=nakano
WorkingDirectory=/home/nakano/PachislotDataCounter/host
# run.sh loads host/.env, launches uvicorn, and mirrors output to run.log.
ExecStart=/bin/bash /home/nakano/PachislotDataCounter/run.sh
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
# App config (PUBLIC_BASE_URL etc.) is read from host/.env — see host/.env.example.
# Leading "-" makes it optional so the service still starts without the file.
EnvironmentFile=-/home/nakano/PachislotDataCounter/host/.env

[Install]
WantedBy=multi-user.target
