# systemd/pachislot-data-counter.service
[Unit]
Description=Pachislot Data Counter (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nakano
Group=nakano
WorkingDirectory=/home/nakano/PachislotDataCounter
ExecStart=/home/nakano/PachislotDataCounter/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
